# config/systemd/user/backup.timer
[Unit]
Description=Daily backup of hard drives

[Timer]
OnCalendar=daily
Persistent=true     
 
[Install]
WantedBy=timers.target
